# systemd/logclean.service
[Unit]
Description=delete logs and networks periodically

[Service]
ExecStartPre=logclean
ExecStart=logclean

[Install]
WantedBy=multi-user.target
